# channels/lxde/airootfs/etc/systemd/system/getty@tty4.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin %USERNAME% --noclear %I $TERM 
