# channels/lxqt/airootfs.ja/etc/skel/.config/autostart/ibus_daemon.desktop
[Desktop Entry]
Exec=ibus-daemon -r -d -x
Icon=system-run
Name=ibus
Path=
Terminal=False
Type=Application
